# artools-workspace/iso-profiles/common/live-overlay/etc/skel/.config/autostart/make-installer-executable.desktop
[Desktop Entry]
Exec=chmod +x "~/Desktop/calamares.desktop"
Name=Make the installer on the desktop executable
StartupNotify=false
Terminal=false
Type=Application
Name[en_US]=make-installer-executable.desktop
