# artools-workspace/iso-profiles/common/live-overlay/etc/local.d/allow_weak_pass_in_calamares.start
:

echo "allowWeakPasswords: true" >> /etc/calamares/modules/users.conf &>/dev/null
echo "allowWeakPasswordsDefault: true" >> /etc/calamares/modules/users.conf &>/dev/null
